# etc/cron.d/blocklist
0 * * * *       root    /bin/bash /etc/fail2ban/blocklist.sh